# deploy/systemd/poe-market-poller.service
[Unit]
Description=PoE Market Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/poe-market-flips
ExecStart=/opt/poe-market-flips/.venv/bin/python -m poller --poll-interval 0
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/poe-market-flips/secrets.env

[Install]
WantedBy=multi-user.target
